# phusion-common/services/php-fpm/run
#!/bin/sh

set -eu

php-fpm${PHP_VERSION}